# askbot/container/crontab
0 * * * * /{{ askbot_site }}/{{ askbot_app }}/cron-askbot.sh
